# deploy/rust-dns.service
[Unit]
Description=rust-dns blocking DNS resolver
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Dedicated unprivileged user that owns /opt/rust-dns (created by install.sh).
User=rust-dns
Group=rust-dns
# Working dir holds config.toml, blocklist.txt, cache.redb and logs/.
WorkingDirectory=/opt/rust-dns
ExecStart=/opt/rust-dns/rust-dns /opt/rust-dns/config.toml
Restart=on-failure
RestartSec=2
Environment=RUST_LOG=info

# Bind port 53 without running as root.
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
NoNewPrivileges=yes

# Sandbox: the filesystem is read-only except /opt/rust-dns, where the server
# writes the generated token, the cache, and query logs.
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ReadWritePaths=/opt/rust-dns
# Files the server writes (token in config.toml, cache, logs) stay owner-only.
UMask=0077

# Graceful stop lets the cache flush.
KillSignal=SIGINT
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
